# debian/isso.socket
[Unit]
Description=isso socket

[Socket]
ListenStream=/run/isso.sock

[Install]
WantedBy=sockets.target
